# scripts/systemd/sip-bridge.service
[Unit]
Description=Ekaette SIP Bridge (Africa's Talking)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mac
WorkingDirectory=/home/mac
Environment=PYTHONPATH=/home/mac
EnvironmentFile=/home/mac/sip_bridge.env
ExecStart=/home/mac/sip_venv/bin/python -m sip_bridge.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
